# ansible/layer2/roles/podman_runtime_workload/templates/nginx-proxy.container.j2
# nginx reverse-proxy — sole published port; routes to Gorobi, ingestion, FastAPI, Dagster.
[Unit]
Description=nginx reverse proxy for data platform runtime stack
After={{ runtime_stack_network_name }}.service network-online.target
Requires={{ runtime_stack_network_name }}.service
{% for svc in runtime_workload_services %}
Wants={{ svc.name }}.service
{% endfor %}
Wants=network-online.target

[Container]
ContainerName=nginx-proxy
Image={{ runtime_nginx_local_image }}
Network={{ runtime_stack_network_name }}
PublishPort={{ runtime_nginx_publish_port }}:8080
LogDriver=journald
Memory={{ runtime_nginx_memory_mb }}m
Cpus={{ runtime_nginx_cpus }}
NoNewPrivileges=true
Volume=/etc/nginx/runtime/nginx.conf:/etc/nginx/nginx.conf:ro,Z
Volume=/etc/nginx/runtime/conf.d:/etc/nginx/conf.d:ro,Z
HealthCmd=wget -qO- http://127.0.0.1:8080/health || exit 1
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=15s
HealthTimeout=5s

[Service]
Restart=on-failure
RestartSec=5s
TimeoutStartSec=90

[Install]
WantedBy=multi-user.target default.target
